# scripts/com.aily.chaos.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aily.chaos</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/luzi/code/aily/scripts/run_chaos_daemon.py</string>
        <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/luzi/code/aily</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/luzi/code/aily</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/luzi/aily_chaos/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/luzi/aily_chaos/daemon.error.log</string>

    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
