# installer/com.aily.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aily.app</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Aily.app/Contents/MacOS/Aily</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/luzi/.aily/aily.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/luzi/.aily/aily.error.log</string>
    <key>WorkingDirectory</key>
    <string>/Applications/Aily.app/Contents/Resources</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
